# docker/supervisord.conf
# ============================================================================
# Supervisor Configuration - Sistema Trading Bitcoin com Ollama LLM
# ============================================================================
# Configuração do supervisor para gerenciar processos no container
# ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/opt/bitcoin-trading-system/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/opt/bitcoin-trading-system/logs/

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ============================================================================
# Ollama Service
# ============================================================================
[program:ollama]
command=ollama serve
user=bitcoin-trader
directory=/opt/bitcoin-trading-system
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/opt/bitcoin-trading-system/logs/ollama.log
stderr_logfile=/opt/bitcoin-trading-system/logs/ollama_error.log
environment=HOME="/home/bitcoin-trader",USER="bitcoin-trader"
priority=100

# ============================================================================
# Trading System Service
# ============================================================================
[program:trading-system]
command=/opt/bitcoin-trading-system/venv/bin/python /opt/bitcoin-trading-system/src/bitcoin_trading_system_with_ollama.py
user=bitcoin-trader
directory=/opt/bitcoin-trading-system
autostart=true
autorestart=true
startretries=3
startsecs=30
stdout_logfile=/opt/bitcoin-trading-system/logs/trading.log
stderr_logfile=/opt/bitcoin-trading-system/logs/trading_error.log
environment=HOME="/home/bitcoin-trader",USER="bitcoin-trader",PATH="/opt/bitcoin-trading-system/venv/bin:%(ENV_PATH)s"
priority=200
depends_on=ollama

# ============================================================================
# System Monitor Service
# ============================================================================
[program:monitor]
command=/opt/bitcoin-trading-system/venv/bin/python /opt/bitcoin-trading-system/src/monitor.py
user=bitcoin-trader
directory=/opt/bitcoin-trading-system
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/opt/bitcoin-trading-system/logs/monitor.log
stderr_logfile=/opt/bitcoin-trading-system/logs/monitor_error.log
environment=HOME="/home/bitcoin-trader",USER="bitcoin-trader",PATH="/opt/bitcoin-trading-system/venv/bin:%(ENV_PATH)s"
priority=300
depends_on=trading-system

# ============================================================================
# Health Check Service
# ============================================================================
[program:healthcheck]
command=/bin/bash -c 'while true; do curl -f http://localhost:11434/api/tags > /dev/null 2>&1 || echo "$(date): Ollama health check failed" >> /opt/bitcoin-trading-system/logs/health.log; sleep 30; done'
user=bitcoin-trader
directory=/opt/bitcoin-trading-system
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/opt/bitcoin-trading-system/logs/health.log
stderr_logfile=/opt/bitcoin-trading-system/logs/health_error.log
priority=400

# ============================================================================
# Log Rotation Service (runs once daily)
# ============================================================================
[program:logrotate]
command=/bin/bash -c 'while true; do sleep 86400; find /opt/bitcoin-trading-system/logs -name "*.log" -size +100M -exec truncate -s 50M {} \;; done'
user=bitcoin-trader
directory=/opt/bitcoin-trading-system
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/opt/bitcoin-trading-system/logs/logrotate.log
stderr_logfile=/opt/bitcoin-trading-system/logs/logrotate_error.log
priority=500

